# devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/home/f1tenthteam2/pleasework/src/ros/rosbash/rosbash"
