# devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/f1tenthteam2/pleasework/src/ros/rosbash/rostcsh"
